# debian/rustguac.service
[Unit]
Description=rustguac web session proxy
After=network.target rustguac-guacd.service
Requires=rustguac-guacd.service

[Service]
Type=simple
User=rustguac
WorkingDirectory=/opt/rustguac
ExecStart=/opt/rustguac/bin/rustguac --config /opt/rustguac/config.toml serve
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
EnvironmentFile=-/opt/rustguac/env

[Install]
WantedBy=multi-user.target
